# home/.config/systemd/user/ssh-agent.service
[Unit]
Description=SSH Key Agent

[Service]
Type=simple
ExecStart=ssh-agent -D -a %t/ssh-agent.sock

[Install]
WantedBy=default.target
